# deploy/linux/agenteum-net.service
[Unit]
Description=Agenteum Net MCP Server
After=network.target

[Service]
Type=simple
User=%USER%
Group=%GROUP%
WorkingDirectory=%PROJECT_DIR%

Environment=PYTHONPATH=%PROJECT_DIR%
Environment=PATH=%VENV_BIN%:/usr/local/bin:/usr/bin:/bin

ExecStart=%VENV_BIN%/agenteum-net

Restart=always
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

StandardOutput=append:%PROJECT_DIR%/logs/agenteum-net.out.log
StandardError=append:%PROJECT_DIR%/logs/agenteum-net.err.log

NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=%PROJECT_DIR%
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
